# com.scout.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scout.dev</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/q/Developer/scout/start-scout.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/q/Developer/scout/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/q/Developer/scout/logs/launchd-stderr.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>ExitTimeOut</key>
    <integer>10</integer>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>AbandonProcessGroup</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/q/Developer/scout/apps/web</string>
</dict>
</plist>
